# usda-vision-camera.service
[Unit]
Description=USDA Vision Camera System
After=network.target
Wants=network.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/home/alireza/USDA-vision-cameras
ExecStart=/bin/bash /home/alireza/USDA-vision-cameras/start_system.sh
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Environment variables
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
Environment=HOME=/root

# Security settings (adjust as needed)
NoNewPrivileges=false
PrivateTmp=false

[Install]
WantedBy=multi-user.target
